# rk3588_sdk/aarch64-linux-gnu/libc/var/lib/dpkg/info/firefly-fan.prerm
#!/bin/sh

set -e

case "$1" in
	upgrade)
		systemctl disable firefly-fan
		;;
	remove)
		systemctl stop firefly-fan
		systemctl disable firefly-fan
		;;
esac


