# rk3588_sdk/aarch64-linux-gnu/libc/var/lib/dpkg/info/firefly-rk3588npu-driver.postinst
#!/bin/sh

set -e

case "$1" in
	configure)
	systemctl enable rknn_server
esac


