# rk3588_sdk/aarch64-linux-gnu/libc/var/lib/dpkg/info/firefly-rk3588npu-driver.prerm
#!/bin/sh

set -e

case "$1" in
	configure)
	systemctl stop rknn_server
	systemctl disable rknn_server
esac


