# rk3588_sdk/aarch64-linux-gnu/libc/var/lib/dpkg/info/firefly-rk3588npu-driver.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rknn_server.service
/usr
/usr/bin
/usr/bin/restart_rknn.sh
/usr/bin/rknn_server
/usr/bin/start_rknn.sh
/usr/include
/usr/include/rknn_api.h
/usr/include/rknn_matmul_api.h
/usr/lib
/usr/lib/librknnrt.so
/usr/lib/librknn_api.so
