# rk3588_sdk/aarch64-linux-gnu/libc/var/lib/dpkg/info/firefly-wdt.prerm
#!/bin/sh

set -e

case "$1" in
	upgrade)
		;;
	remove)
		systemctl disable firefly-wdt
		;;
esac


