# rk3588_sdk/aarch64-linux-gnu/libc/var/lib/dpkg/info/firefly-fan.postinst
#!/bin/sh

set -e

case "$1" in
	configure)
	systemctl enable firefly-fan
	systemctl start firefly-fan
esac


